# deploy/systemd/perfectweb.service
[Unit]
Description=Perfectweb Astro Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/perfectweb
Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=3000
EnvironmentFile=-/var/www/perfectweb/.env.production
ExecStart=/usr/bin/node /var/www/perfectweb/dist/server/entry.mjs
Restart=always
RestartSec=2
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
